# config/release-watch.service
[Unit]
Description=GitHub Node Release Watcher Bot

[Service]
WorkingDirectory=
ExecStart=

[Install]
WantedBy=default.target
